# git/configurations/supervisor/supervisord.conf
[supervisord]
nodaemon = true

[program:sshd]
command = /usr/sbin/sshd -D
autostart = true
autorestart = true

[program:proxy]
command = /usr/sbin/nginx -D
autostart = true
autorestart = true